# ops/systemd/morse-server.user.service
[Unit]
Description=Topology-First CAD morse-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/topology-first-cad
Environment=PORT=8790
ExecStart=/home/orbiter1/projects/topology-first-cad/target/release/morse-server
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
